# Parser_Pi/boot.service
[Unit]
Description=Run my CMake-built C++ Program at Boot
After=network.target

[Service]
ExecStart=/usr/local/bin/my_program
Restart=on-failure
User=YOUR_USERNAME

[Install]
WantedBy=multi-user.target