# low-cost-sensors.service
[Unit]
Description=Low-Cost Sensors Data Collection System
After=network.target multi-user.target

[Service]
Type=forking
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Low-Cost-Sensors
ExecStart=/home/ubuntu/Low-Cost-Sensors/runOnBoot.sh start
ExecStop=/home/ubuntu/Low-Cost-Sensors/runOnBoot.sh stop
ExecReload=/home/ubuntu/Low-Cost-Sensors/runOnBoot.sh restart
Restart=on-failure
RestartSec=5

# Environment variables
Environment=ROS_MASTER_URI=http://localhost:11311
Environment=ROS_HOSTNAME=localhost

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
